# podconnect/rootfs/etc/s6-overlay/s6-rc.d/owntone/run
#!/usr/bin/with-contenv bashio
# Wait up to ~15s for Avahi to be ready on D-Bus (bounded; never hangs)
for i in $(seq 1 30); do
  dbus-send --system --dest=org.freedesktop.Avahi --print-reply / \
    org.freedesktop.Avahi.Server.GetState >/dev/null 2>&1 && break
  sleep 0.5
done
exec owntone -f -c /etc/owntone.conf
