# podconnect/rootfs/etc/s6-overlay/s6-rc.d/owntone/run
#!/usr/bin/with-contenv bashio
exec owntone -f -c /etc/owntone.conf
