# podconnect/rootfs/etc/s6-overlay/s6-rc.d/avahi/run
#!/usr/bin/with-contenv bashio
exec avahi-daemon --no-chroot --no-drop-root
