# podconnect/rootfs/etc/s6-overlay/s6-rc.d/gl-watchdog/run
#!/usr/bin/with-contenv bashio
# go-librespot has known hang bugs (#300/#240) where the process stays alive but the Spotify
# session is stuck/lost — a plain supervisor won't catch that. Poll its /status and restart
# the service if it's unresponsive for ~90s. s6-svc -r restarts only this service (rc state
# stays "up", so this watchdog isn't taken down with it).
URL="http://localhost:3678/status"
sleep 90   # grace period for first start
fails=0
while true; do
  if curl -fsS -m 5 -o /dev/null "$URL" 2>/dev/null; then
    fails=0
  else
    fails=$((fails + 1))
    bashio::log.warning "go-librespot health check failed (${fails}/3)"
    if [ "$fails" -ge 3 ]; then
      bashio::log.warning "go-librespot unresponsive — restarting it"
      s6-svc -r /run/service/go-librespot 2>/dev/null \
        || s6-svc -r /var/run/service/go-librespot 2>/dev/null \
        || bashio::log.error "watchdog: could not find go-librespot service dir"
      fails=0
      sleep 60
    fi
  fi
  sleep 30
done
