# podconnect/rootfs/etc/s6-overlay/s6-rc.d/select-homepod/run
#!/usr/bin/with-contenv bashio
# Selects the HomePod AirPlay output in OwnTone. Forgiving: auto-selects if there's exactly
# one AirPlay device, matches the name case-insensitively, and logs the names it finds so the
# user can copy the exact one. Does NOT touch volume — OwnTone remembers each speaker's last level.
URL="http://localhost:3689"
# The Ingress picker (manager) writes the chosen device to selected_output.json. It wins over the
# static homepod_name config option, which is now just a fallback for users who prefer typing.
read_target() {
  local saved cfg
  saved="$(jq -r '.name // ""' /data/selected_output.json 2>/dev/null)"
  cfg="$(jq -r '.homepod_name // ""' /data/options.json)"
  if [ -n "$saved" ]; then echo "$saved"; else echo "$cfg"; fi
}
listed=""
lasttarget="__unset__"
while true; do
  # Grace-release: while the manager has freed the HomePod for other apps, don't re-grab it.
  if [ -f /data/released ]; then sleep 5; continue; fi
  # Re-read each loop so picking in the Ingress panel takes effect live (no add-on restart).
  HOMEPOD="$(read_target)"
  if [ "$HOMEPOD" != "$lasttarget" ]; then
    if [ -n "$HOMEPOD" ]; then
      bashio::log.info "select-homepod: target '${HOMEPOD}'"
    else
      bashio::log.info "select-homepod: no HomePod chosen yet — open the PodConnect panel to pick one (auto-selects if exactly one AirPlay device is found)"
    fi
    lasttarget="$HOMEPOD"
  fi
  if curl -fsS -o /dev/null "${URL}/api/config" 2>/dev/null; then
    # OwnTone reports HomePods as "AirPlay 2" (not "AirPlay") — match the prefix, not an exact string.
    airplay="$(curl -fsS "${URL}/api/outputs" 2>/dev/null | jq -c '[.outputs[]? | select((.type // "" | ascii_downcase) | startswith("airplay"))]')"
    [ -z "$airplay" ] && airplay="[]"
    count="$(echo "$airplay" | jq 'length')"
    names="$(echo "$airplay" | jq -r '[.[].name] | join(", ")')"
    if [ "$count" -gt 0 ] && [ "$names" != "$listed" ]; then
      bashio::log.info "AirPlay devices found: ${names}"
      listed="$names"
    fi
    row=""
    if [ -n "$HOMEPOD" ]; then
      row="$(echo "$airplay" | jq -c --arg n "$HOMEPOD" 'map(select((.name // "" | ascii_downcase) == ($n | ascii_downcase)))[0] // empty')"
    elif [ "$count" = "1" ]; then
      row="$(echo "$airplay" | jq -c '.[0]')"
    fi
    if [ -n "$row" ]; then
      ID="$(echo "$row" | jq -r '.id')"
      SEL="$(echo "$row" | jq -r '.selected')"
      NM="$(echo "$row" | jq -r '.name')"
      AUTH="$(echo "$row" | jq -r '(.needs_auth_key // .requires_auth) // false')"
      if [ "$AUTH" = "true" ]; then
        bashio::log.warning "HomePod '${NM}' needs AirPlay verification — Apple Home app > this HomePod > 'Allow Speaker & Display Access' > 'Anyone on the Same Network'."
      elif [ "$SEL" != "true" ]; then
        if curl -fsS -X PUT "${URL}/api/outputs/set" --data "{\"outputs\":[\"${ID}\"]}" >/dev/null 2>&1; then
          bashio::log.info "Selected HomePod '${NM}' (volume left at its last setting)."
        fi
      fi
    fi
  fi
  sleep 10
done
