# podconnect/rootfs/etc/s6-overlay/s6-rc.d/go-librespot/run
#!/usr/bin/with-contenv bashio
exec go-librespot -config_dir /data/go-librespot
